# pim/vdirsyncer-sync.service
[Unit]
Description=vdirsyncer sync

[Service]
ExecStart=vdirsyncer sync
